# conf/valhalla.conf
[supervisord]
logfile=/app/tmp_data/logs/supervisor.log
logfile_maxbytes=10MB
logfile_backups=10
pidfile=/var/run/supervisord.pid

[program:build_loop]
command="/usr/local/bin/run_valhalla.sh"
user=root
autostart=false
autorestart=false
redirect_stderr=true
# Either log to file inside the container or
# log to PID 1 (gunicorn in this case) so docker logs will show it
stdout_logfile=/app/tmp_data/logs/builder.log
stdout_logfile_maxbytes=10MB
# stdout_logfile=/proc/1/fd/1
# stdout_logfile_maxbytes=0
environment=no_proxy="%(ENV_NO_PROXY)s","%(ENV_HTTPS_PROXY)s",http_proxy="%(ENV_HTTP_PROXY)s",MAX_CACHE_SIZE="%(ENV_MAX_CACHE_SIZE)s",CONCURRENCY="%(ENV_CONCURRENCY)s",DATA_DIR="/app/data",TMP_DATA_DIR="/app/tmp_data"
